# src/project_templates/virtual_window/virtual-window-check.timer
[Timer]
OnBootSec=30sec
OnUnitActiveSec=10min
Unit=virtual-window-check.service

[Install]
WantedBy=timers.target
